# web2fb.service
[Unit]
Description=web2fb - Web to Framebuffer Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kiosk
WorkingDirectory=/home/kiosk/web2fb
EnvironmentFile=/home/kiosk/web2fb/.env
ExecStart=/usr/bin/node /home/kiosk/web2fb/web2fb.js

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=web2fb

# Security: Run with limited privileges
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
